# Multipass/configs/certbot-auto-renewal.service
[Unit]
Description=Renew certbot Certificates

[Service]
Type=simple
User=certbot
ExecStart=ExecStart=/usr/bin/certbot renew -q